# examples/problems/argumentation/infrastructure/solvers/files/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

* * * * * root /opt/update-results.sh >> /var/log/update-results-cron.log 2>&1

